# cluster-node/etc/systemd/system/docker.service.d/docker-opts.conf
[Service]
ExecStart=
ExecStart=/usr/bin/docker -d --selinux-enabled -H fd:// $OPTIONS 
